# config/99-brotherql.rules
# /etc/udev/rules.d/99-brotherql.rules

ACTION=="add", SUBSYSTEMS=="usb", KERNEL=="lp0", RUN+="/usr/bin/disable_ql-poweroff.sh /dev/usb/lp0"